# zrok-tunnel.service
[Unit]
Description=Zrok Tunnel for WhatsApp Bot
After=network.target whatsapp-bot.service
PartOf=whatsapp-bot.service

[Service]
Type=simple
User=bojrick
WorkingDirectory=/home/bojrick/reeva-whatsapp-crm-erp/site-bot
Environment="HOME=/home/bojrick"
Environment="ZROK_API_ENDPOINT=https://api-v1.zrok.io"
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/zrok share public localhost:3011 --headless
Restart=always
RestartSec=10
StandardOutput=append:/home/bojrick/reeva-whatsapp-crm-erp/site-bot/zrok.log
StandardError=append:/home/bojrick/reeva-whatsapp-crm-erp/site-bot/zrok.log

[Install]
WantedBy=multi-user.target 